# cloudmon/ansible/project/roles/carbonapi/templates/carbonapi-watcher.service.j2
[Unit]
Description=Carbonapi restarter
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart srv.service

[Install]
WantedBy=multi-user.target
