# cloudmon/ansible/project/roles/carbonapi/templates/carbonapi-watcher.path.j2
[Path]
PathModified=/etc/graphite/ssl

[Install]
WantedBy=multi-user.target
